# butterfly.socket
[Socket]
ListenStream=57575

[Install]
WantedBy=sockets.target
